# etc/mori.service
[Unit]
Description=Automatically put osu! related archives in their places.

[Service]
Type=simple
ExecStart=mori

[Install]
WantedBy=default.target
